# deploy/ronitnath-backup.timer
[Unit]
Description=Nightly ronitnath off-host backup

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=1800
Persistent=true
Unit=ronitnath-backup.service

[Install]
WantedBy=timers.target
