# deploy/ronitnath-restore-drill.timer
[Unit]
Description=Weekly ronitnath off-host backup restore drill

[Timer]
OnCalendar=Sun *-*-* 04:30:00
RandomizedDelaySec=1800
Persistent=true
Unit=ronitnath-restore-drill.service

[Install]
WantedBy=timers.target
